# systemd/jarvis.service
[Unit]
Description=Jarvis Local Voice Agent
After=default.target sound.target
Wants=sound.target

[Service]
Type=simple
ExecStart=%h/.local/share/jarvis/venv/bin/python3 -m jarvis.daemon
Restart=on-failure
RestartSec=3
StartLimitIntervalSec=60
StartLimitBurst=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=default.target
